# final/templates/alerts_lateral_movement.json
[
  {
    "alert": "PsExec Remote Execution Detected",
    "response": "PsExec was used to create and execute a service on a remote host, potentially enabling lateral movement across the network.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.002",
    "recommendation": "Review remote service creation activity, validate administrative access, and investigate related lateral movement attempts."
  },
  {
    "alert": "PsExec Remote Execution Detected",
    "response": "PsExec was used to create and execute a service on a remote system through the Service Control Manager, potentially enabling lateral movement across the network.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.002",
    "recommendation": "Review remote service creation activity, validate administrative access, inspect executed commands, and investigate related lateral movement attempts."
  },
  {
    "alert": "SMBExec Remote Command Execution Detected",
    "response": "SMBExec was used to execute commands on a remote host via SMB and the Service Control Manager, indicating potential lateral movement activity.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.002",
    "recommendation": "Review SMB administrative share activity, validate remote command execution, and investigate the source account and host."
  },
  {
    "alert": "PowerShell Remoting Activity Detected",
    "response": "PowerShell Remoting was used to execute commands on a remote system, potentially enabling lateral movement within the environment.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.006",
    "recommendation": "Review PowerShell remoting sessions, verify administrative access, inspect executed commands, and investigate unauthorized remote activity."
  },
  {
    "alert": "WMI Remote Execution Detected",
    "response": "Windows Management Instrumentation (WMI) was used to execute commands on a remote host, a common technique used for lateral movement.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1047",
    "recommendation": "Review WMI execution activity, validate administrative access, inspect executed commands, and investigate the source system."
  },
  {
    "alert": "DCOM Remote Execution Detected",
    "response": "Distributed Component Object Model (DCOM) was used to execute code on a remote system, potentially enabling stealthy lateral movement.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.003",
    "recommendation": "Review DCOM activity, identify the originating account, inspect spawned processes, and investigate unauthorized remote execution."
  },
  {
    "alert": "Remote Scheduled Task Creation Detected",
    "response": "A scheduled task was created to execute commands on a remote system, potentially facilitating lateral movement or remote code execution.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1053.005",
    "recommendation": "Review the created task, inspect configured commands and triggers, verify administrative access, and investigate related remote activity."
  },
  {
    "alert": "SMB Administrative Share File Transfer Detected",
    "response": "A file was copied to a remote host using SMB administrative shares, potentially preparing a payload for remote execution and lateral movement.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1570",
    "recommendation": "Review SMB file transfer activity, validate administrative access, inspect transferred files, and investigate subsequent execution attempts."
  },
  {
    "alert": "Remote Service Creation Detected",
    "response": "A Windows service was created on a remote system through the Service Control Manager, potentially enabling remote code execution and lateral movement.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.002",
    "recommendation": "Review newly created services, inspect service binary paths, validate administrative access, and investigate remote execution activity."
  },
  {
    "alert": "WinRM Remote Management Activity Detected",
    "response": "Windows Remote Management (WinRM) was used to execute commands or manage a remote system, potentially facilitating lateral movement.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1021.006",
    "recommendation": "Review WinRM sessions, verify administrative access, inspect executed commands, and investigate unauthorized remote management activity."
  },
  {
    "alert": "SharpWMI Remote Execution Detected",
    "response": "SharpWMI was used to execute commands on a remote host through Windows Management Instrumentation, indicating possible lateral movement activity.",
    "severity": "High",
    "category": "Lateral Movement",
    "mitre_attack": "T1047",
    "recommendation": "Review WMI execution logs, inspect executed commands, verify administrative access, and investigate the source system."
  }
]